# pypicker/systemd/buzzipick.service
[Unit]
Description=Buzz iPick

[Service]
User=pi
Group=pi
Type=simple
RuntimeDirectory=/home/pi/buzzipick/pypicker
ExecStart=/home/pi/buzzipick/pypicker/venv/bin/python /home/pi/buzzipick/pypicker/buzzimage.py /media/pi/photos

[Install]
WantedBy=multi-user.target